# packages/kernel/drivers/ap6611-bluetooth/ap6611s-bluetooth.service
[Unit]
Description=Bluetooth AP6611S
After=bluetooth.target

[Service]
Type=simple
ExecStartPre=rfkill unblock all
ExecStart=/usr/bin/brcm_patchram_plus --bd_addr_rand --enable_hci --no2bytes --use_baudrate_for_download --tosleep 200000 --baudrate 1500000 --patchram /lib/firmware/ap6611s/SYN43711A0.hcd /dev/ttyS7

[Install]
WantedBy=multi-user.target